# openvpn-killswitch.service
[Unit]
Description=OpenVPN Killswitch Service
Documentation=https://github.com/my7bropx/K.LINUX
After=network-online.target
Wants=network-online.target
Before=network.target

[Service]
Type=notify
NotifyAccess=main

# Execution
ExecStart=/usr/local/bin/openvpn-killswitch.sh start
ExecStop=/usr/local/bin/openvpn-killswitch.sh stop
ExecReload=/usr/local/bin/openvpn-killswitch.sh restart

# Restart behavior
Restart=on-failure
RestartSec=15

# Timeouts - increased to allow VPN to connect
TimeoutStartSec=120
TimeoutStopSec=30
TimeoutSec=0

# Security settings
PrivateTmp=true
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openvpn-killswitch

# Process settings
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
FinalKillSignal=SIGKILL

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target
